# upstart/gridfs-http.conf
description "GridFS HTTP Server"

start on local-filesystems
stop on shutdown
respawn

script
    /usr/bin/gridfs-http-frontend --listen_address=0.0.0.0 --listen_port=3000 --backend=mongodb://gridfs.local:27117/gridfs
end script
